# Ninja_SAC/WOI/scripts/launchd/com.woi.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.woi.healthcheck</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>cd /Users/Santiago/Desktop/Ai/Ninja_SAC/WOI/woi-listener &amp;&amp; node src/scripts/healthcheck.js || curl -X POST -H 'Content-Type: application/json' -d '{"text":":rotating_light: WOI healthcheck externo FAIL"}' "$SLACK_HEALTHCHECK_WEBHOOK"</string>
  </array>

  <!-- Cada 15 minutos -->
  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/Users/Santiago/Desktop/Ai/Ninja_SAC/WOI/logs/healthcheck.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/Santiago/Desktop/Ai/Ninja_SAC/WOI/logs/healthcheck.log</string>
</dict>
</plist>
